# package/steer_debian/postinst
# script executed following the install

echo "running attlrqlite post initialisation..."
if which update-rc.d >/dev/null 
then
	echo "adding tegu rc.d script to system startup via update-rc.d"
	update-rc.d tegu  start 90 2 3 4 5 . stop 10 0 1 6 .
	if [ $? -ne 0 ]
	then
		echo "unable to add tegu rc.d to the system startup"
	fi

	#TODO -- add support for other rc.d management tools
fi

echo "post initialisation completed"
